# scripts/worker.service
[Unit]
Description=Site Monitoring Worker
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/YOUR_USERNAME/mfa/site-monitoring-worker
ExecStart=/usr/bin/node scripts/worker-runner.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=8080

# Memory limit for 1GB VM
MemoryMax=800M
MemoryHigh=700M

[Install]
WantedBy=multi-user.target
